# build/dpkg/DEBIAN/postinst
#!/bin/bash
# postinst for docker-dns

PATH=$PATH:/bin:/usr/bin:/sbin:/usr/sbin
SERVICE_NAME=docker-dns
CONFIG_FILE="/etc/docker-dns/docker-dns.conf"

get_custom_ip() {
    if [ -f "$CONFIG_FILE" ]; then
        custom_ip=$(grep '^IP=' "$CONFIG_FILE" | head -n1 | cut -d'=' -f2)
        echo "$custom_ip"
    else
        echo ""
    fi
}

get_custom_tld() {
    if [ -f "$CONFIG_FILE" ]; then
        custom_tld=$(grep '^TLD=' "$CONFIG_FILE" | head -n1 | cut -d'=' -f2)
        echo "$custom_tld"
    else
        echo ""
    fi
}

case $1 in
configure)
    if [ ! -d /run/systemd/system ]; then
        echo "Failed to detect systemd, aborting installation"
        exit 1
    fi

    systemctl daemon-reload

    setcap CAP_NET_BIND_SERVICE=+eip /usr/bin/${SERVICE_NAME}

    systemctl start "${SERVICE_NAME}.service"
    systemctl enable "${SERVICE_NAME}.service"

    CUSTOM_IP=$(get_custom_ip)
    if [ -z "$CUSTOM_IP" ]; then
        echo "No 'IP' defined in ${CONFIG_FILE}. Aborting DNS configuration."
        exit 1
    fi

    CUSTOM_TLD=$(get_custom_tld)
    CUSTOM_TLD="${CUSTOM_TLD:-docker}"

    # Build routing domains from TLD (comma-separated -> space-separated ~prefixed)
    ROUTING_DOMAINS=""
    IFS=',' read -ra TLD_ARRAY <<< "$CUSTOM_TLD"
    for tld in "${TLD_ARRAY[@]}"; do
        tld=$(echo "$tld" | xargs) # trim whitespace
        if [ -n "$tld" ]; then
            ROUTING_DOMAINS="${ROUTING_DOMAINS} ~${tld}"
        fi
    done
    ROUTING_DOMAINS=$(echo "$ROUTING_DOMAINS" | xargs) # trim leading space

    # Prepend `nameserver $CUSTOM_IP` to /etc/resolv.conf (idempotent).
    prepend_resolv_conf() {
        grep -q "$CUSTOM_IP" /etc/resolv.conf 2>/dev/null && \
            echo "detected an entry of $CUSTOM_IP in /etc/resolv.conf, removing.." && \
            grep -v "$CUSTOM_IP" /etc/resolv.conf > /tmp/resolv.tmp && \
            cat /tmp/resolv.tmp > /etc/resolv.conf && rm -f /tmp/resolv.tmp
        { printf '# Generated by docker-dns\nnameserver %s\n' "$CUSTOM_IP"; cat /etc/resolv.conf; } > /tmp/resolv.tmp \
            && cat /tmp/resolv.tmp > /etc/resolv.conf && rm -f /tmp/resolv.tmp
    }

    if systemctl is-active systemd-resolved.service >/dev/null 2>&1; then
        echo "Detected systemd-resolved running. Extending its configuration.."

        if [ ! -d /etc/systemd/resolved.conf.d ]; then
            echo "Creating /etc/systemd/resolved.conf.d folder.."
            mkdir -p /etc/systemd/resolved.conf.d
        fi
        echo "Creating /etc/systemd/resolved.conf.d/docker-dns.conf.."
        cat <<EOF >/etc/systemd/resolved.conf.d/docker-dns.conf
[Resolve]
DNS=$CUSTOM_IP
Domains=$ROUTING_DOMAINS
EOF
        systemctl restart systemd-resolved.service
        echo "systemd-resolved extended with DNS=$CUSTOM_IP and Domains=$ROUTING_DOMAINS."
    elif systemctl is-active NetworkManager.service >/dev/null 2>&1; then
        # NetworkManager active without systemd-resolved (Debian Desktop).
        # NM rewrites /etc/resolv.conf on every connection event and would clobber
        # a one-shot prepend. Install a dispatcher script so NM re-prepends for us
        # after every up/down/dhcp event, and prepend once immediately so DNS
        # works between now and the next NM event.
        echo "Detected NetworkManager without systemd-resolved. Installing dispatcher.."

        mkdir -p /etc/NetworkManager/dispatcher.d
        cat <<EOF >/etc/NetworkManager/dispatcher.d/docker-dns
#!/bin/bash
# Re-prepend docker-dns nameserver after NetworkManager rewrites /etc/resolv.conf.
# NM invokes dispatchers with "\$1 = interface" and "\$2 = action".
# up/dhcp*-change: NM has just written resolv.conf with connection-provided DNS.
# down: NM rewrites resolv.conf to drop the departing interface's DNS servers.
case "\$2" in
    up|down|dhcp4-change|dhcp6-change)
        if ! grep -q '$CUSTOM_IP' /etc/resolv.conf 2>/dev/null; then
            { printf '# Generated by docker-dns\nnameserver %s\n' '$CUSTOM_IP'
              cat /etc/resolv.conf
            } > /tmp/resolv.docker-dns.tmp \\
            && cat /tmp/resolv.docker-dns.tmp > /etc/resolv.conf \\
            && rm -f /tmp/resolv.docker-dns.tmp
        fi
        ;;
esac
EOF
        chmod 0755 /etc/NetworkManager/dispatcher.d/docker-dns
        chown root:root /etc/NetworkManager/dispatcher.d/docker-dns

        prepend_resolv_conf
        echo "NetworkManager dispatcher installed; /etc/resolv.conf prepended with nameserver $CUSTOM_IP."
    else
        echo "Neither systemd-resolved nor NetworkManager is active. Updating /etc/resolv.conf.."
        prepend_resolv_conf
        echo "/etc/resolv.conf updated with nameserver $CUSTOM_IP."
    fi
    ;;
*)
    echo "Usage: $0 configure"
    exit 1
    ;;
esac
