# build/dpkg/DEBIAN/postrm
#!/bin/bash
# postrm for docker-dns

PATH=$PATH:/bin:/usr/bin:/sbin:/usr/sbin
SERVICE_NAME=docker-dns

case $1 in
remove|purge)
    systemctl stop "${SERVICE_NAME}.service" 2>/dev/null || true
    systemctl disable "${SERVICE_NAME}.service" 2>/dev/null || true
    systemctl daemon-reload

    # systemd-resolved drop-in
    rm -f /etc/systemd/resolved.conf.d/docker-dns.conf
    systemctl restart systemd-resolved 2>/dev/null || true

    # dnsmasq
    rm -f /etc/dnsmasq.d/docker-dns.conf
    systemctl restart dnsmasq 2>/dev/null || true

    # unbound
    rm -f /etc/unbound/unbound.conf.d/docker-dns.conf
    systemctl restart unbound 2>/dev/null || true

    # NetworkManager dnsmasq
    rm -f /etc/NetworkManager/dnsmasq.d/docker-dns.conf
    systemctl restart NetworkManager 2>/dev/null || true

    # Plain /etc/resolv.conf (non-resolved systems like Debian)
    if [ -f /etc/resolv.conf ] && grep -q '# Generated by docker-dns' /etc/resolv.conf 2>/dev/null; then
        awk '/# Generated by docker-dns/{skip=1; next} skip{skip=0; next} 1' /etc/resolv.conf > /tmp/resolv.tmp \
            && cat /tmp/resolv.tmp > /etc/resolv.conf && rm -f /tmp/resolv.tmp
    fi
    ;;
esac
